# images/viewer/supervisor/novncproxy.conf
[program:novncproxy]
# Runs websockify directly to proxy for noVNC installed via apt
command=websockify --web /usr/share/novnc/ 6080 localhost:5900
priority=3 # Start after x11vnc
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
