# images/playwright/supervisord.conf
[program:playwright]
# Using fixed version from Dockerfile/entrypoint
command=npx playwright@__PLAYWRIGHT_VERSION__ run-server --port 3000 --host 0.0.0.0
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=":1",PW_TEST_SCREENSHOT_NO_FONTS_READY="1"

[program:chromium-cdp]
command=chromium --headless=new --no-sandbox --remote-debugging-port=9223 --remote-debugging-address=127.0.0.1 --user-data-dir=/tmp/gbox-chromium --remote-allow-origins=*
priority=15
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cdp-proxy]
command=socat TCP-LISTEN:9222,fork,reuseaddr TCP:127.0.0.1:9223
priority=16
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0